# livescore-pilot.service
[Unit]
Description=Livescore Pilot Gunicorn Service
After=network.target

[Service]
User=root
Group=root
WorkingDirectory=/opt/livescore
Environment="FLASK_SECRET_KEY=$(cat /opt/livescore/secret_key.txt)"
ExecStart=/usr/local/bin/gunicorn -c gunicorn_config.py web_interface:app
Restart=always

[Install]
WantedBy=multi-user.target
